# Lampi/pkg/lampi/etc/supervisor/conf.d/lamp_service.conf
[program:lamp_service]
command=/opt/lampi/lamp_service.py
directory=/opt/lampi
environment=HOME="/opt/lampi"
priority=300
autostart=true
autorestart=true
