# Lampi/pkg/lampi/etc/supervisor/conf.d/lamp_ui.conf
[program:lamp_ui]
command=/opt/lampi/lamp_ui.py
directory=/opt/lampi
environment=HOME="/opt/lampi"
user=lampi
priority=400
autostart=true
autorestart=true
